# deploy/smartroom-detect.service
[Unit]
Description=Smartroom person-detection over saved recordings (oneshot)

[Service]
Type=oneshot
WorkingDirectory=/home/alawn/Code/smartroom-control
EnvironmentFile=-/home/alawn/Code/smartroom-control/.env.local
# Long videos + the ffmpeg annotation pass can take a while; don't let systemd kill a legit run.
TimeoutStartSec=infinity
ExecStart=/home/alawn/Code/yolo-bench/.venv/bin/python /home/alawn/Code/smartroom-control/detect/detect.py

# Triggered by the .timer and .path units (and the dashboard) — no [Install].
